# etc/periodic/daily/sitemgr-cleanup
#!/bin/sh
# Daily cleanup of site temporary files

/usr/local/bin/sitemgr-cleanup